# setup/local.inboard-agent.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>local.inboard-agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__INBOARD_HOME__/engines/inbox-agent.sh</string>
  </array>
  <key>StartInterval</key><integer>__PULL_INTERVAL__</integer>
  <key>RunAtLoad</key><false/>
  <key>ProcessType</key><string>Background</string>
  <key>StandardOutPath</key><string>__INBOARD_HOME__/logs/launchd.out</string>
  <key>StandardErrorPath</key><string>__INBOARD_HOME__/logs/launchd.err</string>
</dict>
</plist>
